# roles/utils/templates/systemd.service.j2
[Unit]
Description={{name}}

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart={{command}}

[Install]
WantedBy=graphical-session.target
